# redis.service
[Unit]
Description=Redis

[Service]
TimeoutStartSec=0
ExecStart=/usr/redis/src/redis-server