# vagrant_data/grafana.service
[Unit]
Description=grafana

[Service]
TimeoutStartSec=0
ExecStart=/usr/grafana/bin/grafana-server -config=/usr/grafana/conf/custom.ini -homepath=/usr/grafana web


[Install]
RequiredBy=grafana.service