# vagrant_data/heka.service
[Unit]
Description=Heka
After=syslog-ng.service
Requires=syslog-ng.service

[Service]
TimeoutStartSec=0
ExecStart=/usr/heka/bin/hekad

[Install]