# vagrant_data/node-exporter.service
[Unit]
Description=node_exporter

[Service]
TimeoutStartSec=0
ExecStart=/usr/node_exporter/node_exporter -web.listen-address 'localhost:9100' 

[Install]